# packages/dev-python/vdirsyncer/files/systemd_user/vdirsyncer.timer
[Unit]
Description=Sync CardDAV and CalDAV accounts (hourly)

[Timer]
OnActiveSec=1s
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

